# crontab
10 1 * * * cd /app && flock -n /tmp/procnum_daily.lock -c 'uv run python src/main.py'